# server/cloudflared-tunnel.service
[Unit]
Description=Cloudflare Tunnel for Resort Apps
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/cloudflared --config /home/user/.cloudflared/config.yml tunnel run prototype-resort-apps
Restart=always
RestartSec=5
User=cloudflared
Group=cloudflared
Environment=CLOUDFLARED_TOKEN=
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target